# .config/systemd/user/clash.service
[Unit]
Description=A rule-based tunnel in Go.

[Service]
ExecStart=/usr/bin/clash

[Install]
WantedBy=default.target
