# .config/systemd/user/syncthing.service
[Unit]
Description=A continuous file synchronization program.

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
